# deployment/supervisor/iwy_frontend.conf
[program:iwy-frontend]
command=/bin/bash -c 'npm run dev'
autostart=true
autorestart=true
user=ashymanskyi
directory=/opt/iwy-frontend
environment=HOME="/home/ashymanskyi",USER="ashymanskyi"
